# playbooks/roles/dallenjs/templates/dallenjs.service.j2
[Unit]
Description=DallenJS
Requires=After=mysql.service

[Service]
ExecStart=/srv/app/.nvm/versions/node/v10.16.0/bin/node /srv/app/DallenJS/server/app.js
WorkingDirectory=/srv/app/DallenJS/server
Restart=always
RestartSec=10
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=dallenjs
Environment=NODE_ENV=production
User=app
Group=app

[Install]
WantedBy=multi-user.target